# buildroot/board/canrec/pi4/rootfs_overlay/etc/inittab
# /etc/inittab — BusyBox init
# Root is mounted READ-ONLY from kernel cmdline (ro).
# DO NOT add a "remount,rw /" line here.

# Mount virtual filesystems
::sysinit:/bin/mount -t proc proc /proc
::sysinit:/bin/mount -t sysfs sysfs /sys
::sysinit:/bin/mount -t devtmpfs devtmpfs /dev
::sysinit:/bin/mkdir -p /dev/pts /dev/shm
# Mount /tmp and /run as tmpfs (defined in /etc/fstab)
::sysinit:/bin/mount -a
::sysinit:/bin/hostname -F /etc/hostname

# Run init.d scripts (network, syslog, etc.)
::sysinit:/etc/init.d/rcS

# Seed writable webroot from read-only overlay into tmpfs so rsync can update
# files live without reflashing (/tmp is already a 64 MB tmpfs from fstab)
::sysinit:/bin/cp -a /var/www /tmp/www

# Start the web status server (daemonises itself)
::once:/usr/sbin/httpd -p 80 -h /tmp/www -c /etc/httpd.conf

# HDMI console
tty1::respawn:/sbin/getty -L tty1 0 vt100
# Serial console (UART)
ttyAMA0::respawn:/sbin/getty -L ttyAMA0 115200 vt100

# Clean shutdown
::shutdown:/etc/init.d/rcK
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
::ctrlaltdel:/sbin/reboot
